# deploy/systemd/compliancepulse-worker.service
[Unit]
Description=CompliancePulse Scan Worker
After=network.target
Requires=network.target

[Service]
Type=simple
WorkingDirectory=/opt/compliancepulse/backend
Environment=PYTHONPATH=/opt/compliancepulse/backend
Environment=SESSION_SECRET=changeme-super-secret
ExecStart=/usr/bin/python3 /opt/compliancepulse/backend/worker.py
Restart=always
RestartSec=5
User=compliance
Group=compliance

[Install]
WantedBy=multi-user.target
